# examples/irc/ansible/roles/miniircd/templates/miniircd.service.j2
[Unit]
Description=miniircd

[Service]
ExecStart=/opt/miniircd/miniircd --log-file=/opt/miniircd/logs/miniircd.log --debug
User=miniircd

[Install]
WantedBy=multi-user.target
